# build/relay/com.internetmerge.relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.internetmerge.relay</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/internetmerge</string>
    <string>relay</string>
    <string>--listen</string>
    <string>:7000</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>INTERNETMERGE_RELAY_KEY</key>
    <string>__RELAY_KEY__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/internetmerge-relay.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/internetmerge-relay.log</string>
</dict>
</plist>
